# build/lib/csiplugin-connector.service
[Unit]
Description=csiplugin connector
After=network.target remote-fs.target nss-lookup.target

[Service]
Type=simple
RuntimeDirectory=csi-tool/connector
ExecStart=/etc/csi-tool/csiplugin-connector
Restart=always
RestartSec=5s

[Install]
WantedBy=multi-user.target