# build/lib/csiplugin-connector.service
[Unit]
Description=csiplugin connector
After=network.target remote-fs.target nss-lookup.target

[Service]
Type=forking
RuntimeDirectory=csiplugin
ExecStart=/etc/csi-tool/csiplugin-connector
ExecReload=/bin/kill -s HUP $MAINPID
ExecStop=/bin/kill -s QUIT $MAINPID
Restart=always
RestartSec=5s

[Install]
WantedBy=multi-user.target